# prima.timer
[Unit]
Description=PRIMA — Heartbeat timer for Azure PIM role activation

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
