# prima.timer
[Unit]
Description=PRIMA — Heartbeat timer for Azure PIM role activation

[Timer]
OnActiveSec=2min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
